# docker/tweetheat/dev/etc/service/server/run
#!/usr/bin/env bash

. /etc/environment
exec dev.work-on-backend